# .config/systemd/user/weztrunk-repos-pull.service
[Unit]
Description=WezTrunk conservative repo pull

[Service]
Type=oneshot
ExecStart=%h/.local/bin/weztrunk repos pull --notify
Nice=10
TimeoutStartSec=120
